# debian_conf/normal/etc/profile.d/nobeep.sh
#!/usr/bin/env bash

xset -b                   &> /dev/null

setterm  -blength 0       &> /dev/null
setterm --blength --store &> /dev/null
